# scripts/wa-bridge.service
[Unit]
Description=zcoms WhatsApp bridge (Baileys sidecar)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Edit this path if your zcoms checkout lives elsewhere.
WorkingDirectory=%h/personal/zcoms/whatsapp-bridge
ExecStart=/usr/bin/node index.js
# Socket + auth dir default to ~/.config/zcoms/{wa.sock,wa-auth}; override here if needed:
# Environment=WA_SOCK=%h/.config/zcoms/wa.sock
# Environment=WA_AUTH=%h/.config/zcoms/wa-auth
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
