# scripts/zcoms-daemon.service
[Unit]
Description=zcoms agent bridge (zc init agent)
After=network-online.target
Wants=network-online.target
# If the optional WhatsApp sidecar is enabled, start it first so its socket
# exists before triage runs. Harmless when wa-bridge.service isn't installed.
After=wa-bridge.service
Wants=wa-bridge.service

[Service]
Type=simple
# .env (TG_API_ID/TG_API_HASH) is loaded from the working directory.
# Edit this path if your zcoms checkout lives elsewhere.
WorkingDirectory=%h/personal/zcoms
ExecStart=%h/.local/bin/zc init agent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
